# crontab
* * * * * php /var/www/html/nhc/artisan schedule:run
